# server/pdfdraw.service
[Unit]
Description=PDF Annotation service for Nextcloud
After=network.target

[Service]
WorkingDirectory=/path/to/nextcloud/apps/pdfdraw/server
User=root
ExecStart=/usr/bin/nodejs server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
